# crontab/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (- installed on Tue Dec  4 06:47:47 2018)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
@weekly /usr/local/sbin/dehydrated-renew.sh
